# kvm-no-cloud-init/deb/lib/systemd/system/kvm-no-cloud-init.service
[Unit]
Before=systemd-networkd.service

[Service]
ExecStart=/sbin/kvm-no-cloud-init
Type=oneshot

[Install]
WantedBy=systemd-networkd.service
